# minecraft/optVersion/svc/minecraft-bedrock-backup.service
[Unit]
Description=Backup Minecraft Bedrock World using /root/Backup.sh

[Service]
Type=simple
ExecStart=/bin/bash /root/Backup.sh


