# minecraft/optVersion/svc/siteforest-update.service
[Unit]
Description=Sync the Siteforest

[Service]
Type=simple
ExecStart=/bin/bash /root/updateSiteForest.sh
